# roles/install-app/templates/app-slack.service
[Unit]
Description=APP Slack-O

[Service]
ExecStart={{ app_path }}/app/slacko-api.sh

[Install]
WantedBy=multi-user.target